# configs/cart.service
[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis.devopsb71.online
Environment=CATALOGUE_HOST=catalogue.devopsb71.online
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target
